# recipes-setting/automount/files/automount.service
[Unit]
Description=automount Service

[Service]
ExecStart=/bin/sh -c "/etc/automount/automount.sh"

[Install]
WantedBy=multi-user.target
